# root/etc/s6-overlay/s6-rc.d/init-couchpotato-config/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

# Create data directories (matches original linuxserver.io layout)
mkdir -p /config/data/database
chown -R abc:abc /config
