# root/etc/s6-overlay/s6-rc.d/svc-couchpotato/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

exec \
  s6-setuidgid abc \
  python3 /app/couchpotato/CouchPotato.py \
    --data_dir=/config \
    --console_log \
    --daemon=0
